# qr3_8.conf
[program:qr3_8]
command=python3 /home/curator3/qr3_8/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_8/test.err.log
stdout_logfile=/home/curator3/qr3_8/test.err.log
